# scripts/systemd/vibeguard-gc.service
[Unit]
Description=VibeGuard Scheduled GC
Documentation=https://github.com/majiayu000/vibeguard

[Service]
Type=oneshot
ExecStart=/bin/bash __VIBEGUARD_DIR__/scripts/gc-scheduled.sh
StandardOutput=append:__HOME__/.vibeguard/gc-systemd.log
StandardError=append:__HOME__/.vibeguard/gc-systemd.log
