# scripts/systemd/vibeguard-gc.service
[Unit]
Description=VibeGuard Scheduled GC
Documentation=https://github.com/majiayu000/vibeguard

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'mkdir -p "__HOME__/.vibeguard" && /bin/bash "__VIBEGUARD_DIR__/scripts/gc-scheduled.sh" >> "__HOME__/.vibeguard/gc-systemd.log" 2>&1'
