# scripts/setup/com.vibeguard.health-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vibeguard.health-report</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__VIBEGUARD_DIR__/scripts/health-report-scheduled.sh</string>
        <string>--scheduled</string>
        <string>--days</string>
        <string>__DAYS__</string>
        <string>--scope</string>
        <string>__SCOPE__</string>
        <string>--output-dir</string>
        <string>__OUTPUT_DIR__</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>__HOME__/.vibeguard/health-report-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.vibeguard/health-report-launchd.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
